# deploy/systemd/bigv-twins-blogger.service
[Unit]
Description=BigV-twins BLOGGER MCP server (search / persona / recent / post / list)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dtl/projects/BigV-twins
ExecStartPre=-/usr/bin/fuser -k 8770/tcp
ExecStart=/home/dtl/miniconda3/envs/bigv-twins/bin/python -m bigv_twins.blogger_server
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=append:/home/dtl/projects/BigV-twins/logs/mcp_blogger.log
StandardError=append:/home/dtl/projects/BigV-twins/logs/mcp_blogger.log

[Install]
WantedBy=default.target
