# systemd/bigv-twins-server.service
[Unit]
Description=BigV-twins MCP server (Zhihu blogger retrieval, streamable-http on 127.0.0.1:8770)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dtl/projects/BigV-twins
ExecStart=/home/dtl/miniconda3/envs/bigv-twins/bin/python -m bigv_twins.server
Restart=on-failure
RestartSec=5s
StandardOutput=append:/home/dtl/projects/BigV-twins/logs/mcp_server.log
StandardError=append:/home/dtl/projects/BigV-twins/logs/mcp_server.log

[Install]
WantedBy=default.target
